# private_dot_config/systemd/user/theme-monitor.service
[Unit]
Description=Sync alacritty + zellij themes with GNOME color-scheme
After=graphical-session.target
PartOf=graphical-session.target
ConditionHost=mason-xps

[Service]
Type=simple
ExecStart=%h/.local/bin/theme-monitor.sh
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
